# ubuntu/supervisor.d/ssh-kk_ol.ini
[program:ssh-kk_ol]
command=ssh -C -f -N -g -L 9888:127.0.0.1:9888 jump
stdout_logfile=/Users/haojunyu/.supervisord_log/ssh-kk_ol.log
autostart=true   
autorestart=true 
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
